# systemd/neon-logind-cleanup.timer
[Unit]
Description=Cleanup lost logind sessions

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
