# systemd/neon-logind-cleanup.service
[Unit]
Description=Cleanup lost logind sessions

[Service]
ExecStart=/home/neonarchives/bin/neon-logind-cleanup
WorkingDirectory=/home/neonarchives
